# systemd/linux-netmon-agent.service
[Unit]
Description=Linux Network Monitoring Agent (Collector)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=eli
Group=eli
WorkingDirectory=/home/eli/projects/linux-network-monitoring/agent

Environment=PYTHONUNBUFFERED=1
ExecStart=/home/eli/projects/linux-network-monitoring/agent/.venv/bin/python /home/eli/projects/linux-network-monitoring/agent/agent.py


Environment=NETMON_API_BASE=http://127.0.0.1:8000
Environment=NETMON_API_KEY=CHANGE_ME

Restart=on-failure
RestartSec=5

# Security Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target

